# recipes-bsp/mp-decision/mp-decision/mpdecision.conf
# upstart script for mpdecision

description "Start mpdecision"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
   sleep 120
   /usr/bin/mpdecision
end script


